# Recipea/crontab
# Reset database every Sunday at 4:00 AM UTC
0 4 * * 0 rm -f /app/data/recipea.db